# raspberry_pi/systemd/ai-alert.service
[Unit]
Description=AI Alert Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/raspberry_pi
ExecStart=/usr/bin/python3 /home/pi/raspberry_pi/alerts/alert_engine.py
Restart=always

[Install]
WantedBy=multi-user.target
